# root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash

/usr/sbin/dnsmasq -k -u abc -g abc --no-hosts --no-resolv -C /config/dnsmasq.conf --log-facility=-
